# files/dotfiles/.config/autostart/shepherd.desktop
[Desktop Entry]
Type=Application
Name=Shepherd
Exec=shepherd