# files/dotfiles/.config/autostart/signal.desktop
[Desktop Entry]
Type=Application
Name=Signal
Exec=signal-desktop